# docs/Upstart/hyperion-deciders.conf
description "Hyperion Workflow Decider Pool"

start on runlevel [2345]
stop on runlevel [!2345]

pre-start script
  for i in `seq 1 2`
  do
    start hyperion-decider N=$i
  done
end script
