# docs/Upstart/hyperion-workers.conf
description "Hyperion Workflow Worker Pool"

start on runlevel [2345]
stop on runlevel [!2345]

pre-start script
  for i in `seq 1 5`
  do
    start hyperion-worker N=$i
  done
end script
